# projects/vf-ops-watch/deploy/vf-ops-watch.service
[Unit]
Description=VF Ops Watch health runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/vector-forge/projects/vf-ops-watch
Environment="PATH=/home/ubuntu/.npm-global/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/ubuntu/vector-forge/projects/vf-ops-watch/runner.py --config /home/ubuntu/vector-forge/projects/vf-ops-watch/config.yaml
SuccessExitStatus=1 2

# Hardening (safe defaults for local script)
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/ubuntu/vector-forge/projects/vf-ops-watch/out
